# napt/systemd/napt.service
[Unit]
Description=Napt Network Rule Processor
Documentation=https://github.com/napt/napt
After=network-online.target
Wants=network-online.target

[Service]
User=root
Type=simple
ExecStart=/usr/local/bin/napt -c /var/lib/napt/config.toml
Restart=on-failure
RestartSec=5s
TimeoutStartSec=30s
TimeoutStopSec=30s

# Security hardening
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/lib/napt /var/run
CapabilityBoundingSet=CAP_NET_ADMIN CAP_NET_RAW
AmbientCapabilities=CAP_NET_ADMIN CAP_NET_RAW

# Resource limits
LimitNOFILE=65536
LimitNPROC=102400

# Logging
StandardOutput=journal
StandardError=journal
SyslogIdentifier=Napt

[Install]
WantedBy=multi-user.target
